# services/supervisor/crontab
* * * * * php /www/aggreagtion-api/artisan schedule:run  >> /www/aggreagtion-api/crontab.log  2>&1