# services/supervisor/queue.ini
[program:queue]
command=php /www/aggreagtion-api/artisan queue:work   redis  --tries=3  --max-time=3600  
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/logs/supervisor.log
stopwaitsecs=3600
numprocs=8
process_name=%(program_name)s_%(process_num)02d    